# deploy/salon-mis.service
[Unit]
Description=Salon MIS Gunicorn
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/salon-mis
EnvironmentFile=/var/www/salon-mis/.env
RuntimeDirectory=salon-mis
RuntimeDirectoryMode=0755
ExecStart=/var/www/salon-mis/.venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/run/salon-mis/salon-mis.sock \
    --access-logfile /var/log/salon-mis/access.log \
    --error-logfile /var/log/salon-mis/error.log \
    salon_mis.wsgi:application
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
